# p0sixspwn.tproj/Root/Library/LaunchDaemons/com.apple.MobileFileIntegrity.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>DYLD_INSERT_LIBRARIES</key>
		<string>/Developer/Library/Helpers/_.dylib</string>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.MobileFileIntegrity</key>
		<dict>
			<key>HostSpecialPort</key>
			<integer>18</integer>
		</dict>
	</dict>
	<key>KeepAlive</key>
	<true/>
	<key>Label</key>
	<string>com.apple.MobileFileIntegrity</string>
	<key>POSIXSpawnType</key>
	<string>Interactive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/amfid</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>UserName</key>
	<string>root</string>
</dict>
</plist>
