# p0sixspwn.tproj/Root/Library/LaunchDaemons/com.apple.First.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>KeepAlive</key>
	<true/>
	<key>Label</key>
	<string>com.apple.First</string>
	<key>POSIXSpawnType</key>
	<string>Interactive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/launchctl</string>
		<string>unload</string>
		<string>/System/Library/LaunchDaemons/com.apple.MobileFileIntegrity.plist</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>UserName</key>
	<string>root</string>
	<key>LaunchOnlyOnce</key>
	<true/>
</dict>
</plist>
